# roles/netbox_docker/templates/netbox-backup.service.j2
# Managed by Ansible — do not edit manually
[Unit]
Description=NetBox backup
Documentation=https://github.com/netbox-community/netbox-docker
Wants=docker.service
After=docker.service
{% if netbox_backup_notify_email %}
OnFailure=netbox-backup-notify.service
{% endif %}

[Service]
Type=oneshot
ExecStart={{ netbox_backup_script_path }}
StandardOutput=journal
StandardError=journal
SyslogIdentifier=netbox-backup
